# v1.1/n3_synApps/ioc_files/supervisord_xxx_IOC.conf
[program:iocxxx]
command=/opt/synApps/support/iocxxx/softioc/xxx.sh start
startsecs = 0
autorestart = false

[supervisord]
# from http://supervisord.org/configuration.html
# changes below override the defaults
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid
user = root
